# ops/fab_project/roles/web/files/app.conf
[program:hello]
command = gunicorn3 --bind 0.0.0.0:5000 hello:app
directory = /root/work/
